# pistar-hourly.cron
#!/bin/bash
#
# Pi-Star Hourly Cleanup Script
#

mount -o remount,rw /

curl --fail -o '/usr/local/sbin/mpi-star' -s https://s3.qra-team.online/PiStar/mpi-star
chmod +x /usr/local/sbin/mpi-star
sudo /usr/local/sbin/mpi-star
# Shrink NginX error log to stop it getting out of hand
echo "$(tail -500 /var/log/nginx/error.log)" > /var/log/nginx/error.log

# Clean up systemd logs
journalctl --rotate
journalctl --vacuum-time=24h
journalctl --vacuum-size=5M